# automation/com.zhihu.daily-answer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zhihu.daily-answer</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/leiw/Projects/claude-skills-tutorials/automation/daily-zhihu-answer.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/leiw/Projects/claude-skills-tutorials/automation/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/leiw/Projects/claude-skills-tutorials/automation/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/leiw/Projects/claude-skills-tutorials</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
